# docker-desktop-fix.service.template
[Unit]
Description=Remove restrictions to unpriviledged namespaces for Docker Desktop

[Service]
Type=oneshot
ExecStart=/bin/sh PWD/script.sh

[Install]
WantedBy=multi-user.target
